# docker/rc.local
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

touch /var/lock/subsys/local
mosquitto -d
sleep 2
python /home/ec2-user/Group02_DIT029_H17/docker/docker_mqtt_client.py